# releng/airootfs/etc/skel/.config/autostart/homesymlinks.desktop
[Desktop Entry]
Exec=sh -c "~/.config/autostart/homesymlinks.sh"
Name=Fix installer symlinks
StartupNotify=false
Terminal=false
Type=Application
